# system/modular-monitoring/systemd/modular-monitor.service
[Unit]
Description=Modular System Monitor
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash __SCRIPT_DIR__/orchestrator.sh
User=root
Group=root
WorkingDirectory=__SCRIPT_DIR__

# Security settings
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/tmp /var/log /sys /proc

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=modular-monitor

[Install]
WantedBy=multi-user.target
