# system/debug-watchdogs/critical-monitor.service
[Unit]
Description=Critical System Monitor - Real-time thermal, USB, and memory monitoring
Documentation=man:systemd.service(5)
After=multi-user.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=SCRIPT_PATH_PLACEHOLDER/critical-monitor.sh
User=root
Group=root

# Resource limits (keep it lightweight)
MemoryMax=50M
CPUQuota=5%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=critical-monitor

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/tmp/critical-monitor-state /var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
