# system/debug-watchdogs/critical-monitor.timer
[Unit]
Description=Critical System Monitor Timer - Runs every 1 minute
Documentation=man:systemd.timer(5)
Requires=critical-monitor.service

[Timer]
# Run every 1 minute for faster thermal response
OnCalendar=*:0/1
# Run 30 seconds after boot
OnBootSec=30sec
# If we miss a run, catch up immediately
Persistent=true
# Reduce randomization for faster response
RandomizedDelaySec=5sec

[Install]
WantedBy=timers.target
